# docker/supervisord/namenode.conf
[program:namenode]
command=${HADOOP_HOME}/bin/hdfs namenode ;

[program:resourcemanager]
command=${HADOOP_HOME}/bin/yarn resourcemanager ;
